# jupyterport/roles/jupyterport/files/nvidia-smi-exporter.service
[Unit]
Description=Nvidia SMI Exporter

[Service]
ExecStart=/usr/local/bin/nvidia_smi_exporter 9999

[Install]
WantedBy=multi-user.target
